# deploy/systemd/0xcacti-website.service
# deploy/systemd/0xcacti-website.service
[Unit]
Description=0xcacti Website
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/path/to/your/app
ExecStart=/usr/local/bin/sbcl --load "scripts/develop.lisp"
Restart=always
RestartSec=5

# Security
ProtectSystem=strict
ProtectHome=read-only
NoNewPrivileges=yes

# Resource limits
MemoryMax=1G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
